# host-agent/palworld-host-agent.path
[Unit]
Description=Compatibility watcher for Palworld host agent updates

[Path]
PathChanged=__HOST_AGENT_SOURCE__
Unit=palworld-host-agent.service

[Install]
WantedBy=multi-user.target
