# host-agent/palworld-host-agent.path
[Unit]
Description=Restart Palworld host agent when its bundled code changes

[Path]
PathChanged=__HOST_AGENT_SOURCE__
Unit=palworld-host-agent.service

[Install]
WantedBy=multi-user.target
